# supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /home/docker/code/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx

[program:app-django-q]
command = /usr/bin/python3.5 /home/docker/code/manage.py qcluster
